# customer-database-simulator/services/zmq-to-meilisearch.service
[Unit]
Description=ZeroMQ to Meilisearch Sync
After=network.target customer-simulator.service meilisearch.service

[Service]
User=your-user
WorkingDirectory=/path/to/your/customer-database-simulator
Environment="PATH=/path/to/your/customer-database-simulator/.venv/bin"
ExecStart=/path/to/your/customer-database-simulator/.venv/bin/python zmq_to_meilisearch.py
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
